# jlogrotate.cron
0 1 * * * /usr/local/bin/jlogrotate.sh -j -v -k 180 /var/log/jboss/server.log
