# jlogrotate.cron
0 1 * * * /usr/local/bin/jlogrotate.sh /var/log/jboss/server.log
